# data-sources/pubmed/crontab
00 06 * * * /usr/local/bin/python /usr/src/app/pubmed_crawler.py >> /usr/src/app/cron.log 2>&1
