# data-sources/arxiv/crontab
* * * * * /usr/local/bin/python /usr/src/app/arxiv_crawler.py >> /usr/src/app/cron.log 2>&1
