# telegram_update/crontab
00 10 * * * /usr/local/bin/python /usr/src/app/updater.py >> /usr/src/app/cron.log 2>&1
